# docker/gvm/root/etc/services.d/postgres/run
#!/usr/bin/with-contenv bash

echo "Starting PostgreSQL..."
s6-setuidgid postgres /usr/lib/postgresql/$PG_MAJOR/bin/postgres