# config/update-openvas.cron
1 1 * * * /bin/bash /defaults/update-openvas.sh > /dev/stdout
