# docker/gvm/root/etc/services.d/gsa/run
#!/usr/bin/with-contenv bash

echo "Starting Greenbone Security Assistant..."
exec s6-setuidgid abc gsad --foreground --verbose --listen=$GSA_HOST --port=$GSA_PORT