# docker/base/root/defaults/update-openvas.cron
0 0 * * * /bin/bash /defaults/update-openvas.sh
